# files/default/init/deb-webapp.conf
################################################################################
#
# WARNING!
#
# This file is automatically generated by Chef. for <%= @node['hostname'] %>
# Any local modifications will be overwritten on next chef-client run
#
################################################################################

start on runlevel [2345]
stop on runlevel [!2345]

setgid webapp
setuid webapp
exec /opt/webapp/latest/webapp
